# roles/firth_sftp_docker/templates/sftp-sync.service.j2
# {{ ansible_managed }}
[Unit]
Description=Sync SFTP users from database to users.conf
After=network.target mariadb.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sftp-sync.sh
