# roles/firth_sftp_docker/templates/sftp-sync.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Run SFTP user sync every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=10s

[Install]
WantedBy=timers.target
